# debian/sv/log/run
#!/bin/sh
set -e

LOG=/var/log/dyntftpd

# Ensure that the log directory exists.
test -d "$LOG" || mkdir -p -m2750 "$LOG"

# Execute svlogd
exec svlogd -tt "$LOG"
